# xdg-autostart/autostart/ClipIt.desktop
[Desktop Entry]
Exec=qlipper
Name=Qlipper
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
