# xdg-autostart/autostart/DockX.desktop
[Desktop Entry]
Exec=dockx
Hidden=true
Name=DockX
Type=Application
Version=1.0
