# xdg-autostart/autostart/DockX.desktop
[Desktop Entry]
Exec=dockx
Name=DockX
Type=Application
Version=1.0
Hidden=true

