# xdg-autostart/autostart/Indicador de KDE Connect.desktop
[Desktop Entry]
Exec=kdeconnect-indicator
Name=Indicador de KDE Connect
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
