# xdg-autostart/autostart/KSuperKey.desktop
[Desktop Entry]
Exec=ksuperkey &
Hidden=true
Name=KSuperKey
Type=Application
Version=1.0
X-GNOME-Autostart-enabled=false
Comment=
RunHook=0

