# xdg-autostart/autostart/Load urxvt config.desktop
[Desktop Entry]
Exec=xrdb -merge ~/.Xresources
Name=Load urxvt config
Type=Application
Version=1.0
