# xdg-autostart/autostart/Load wallpaper.desktop
[Desktop Entry]
Exec=~/.fehbg
Name=Load wallpaper
Type=Application
Version=1.0
