# xdg-autostart/autostart/Nemo Desktop.desktop
[Desktop Entry]
Exec=nemo-desktop  &
Hidden=false
Name=Nemo Desktop
Type=Application
Version=1.0
Comment=
RunHook=0

