# xdg-autostart/autostart/Redshift.desktop
[Desktop Entry]
Exec=redshift-gtk
Name=Redshift
Type=Application
Version=1.0
