# xdg-autostart/autostart/Redshift.desktop
[Desktop Entry]
Exec=redshift-gtk
Hidden=true
Name=Redshift
Type=Application
Version=1.0
