# xdg-autostart/autostart/tint2.desktop
[Desktop Entry]
Exec=tint2
Hidden=true
Name=tint2
Type=Application
Version=1.0
