# xdg-autostart/autostart/UDiskie.desktop
[Desktop Entry]
Exec=udiskie
Name=UDiskie
Type=Application
Version=1.0
