# xdg-autostart/autostart/UDiskie.desktop
[Desktop Entry]
Exec=udiskie &
Hidden=true
Name=UDiskie
Type=Application
Version=1.0
Comment=
RunHook=0

