# xdg-autostart/autostart/xfce-polkit.desktop
[Desktop Entry]
Comment=Policykit Authentication Agent
Exec=/usr/libexec/xfce-polkit
Hidden=true
Icon=gtk-dialog-authentication
Name=XFCE PolKit
NotShowIn=GNOME;KDE;
Type=Application

